# etc/skel/.config/autostart/keyboardsetting.desktop
[Desktop Entry]
Name=keyboard setting
Comment=Key Board Setting
Exec=keyboardsetting
Terminal=false
Type=Application
